# offload_service.service
[Unit]
Description=Ethernet Sync GUI Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/ethernet_sync_gui.py
Environment=DISPLAY=:0  # Required for GUI applications
Restart=always

[Install]
WantedBy=multi-user.target